# 99-usb-device.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/local/bin/usb-signal.sh"
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/usr/local/bin/usb-signal.sh"
